# infra/systemd/backend-staging.service.example
# Example systemd unit for backend (staging).
# Copy to /etc/systemd/system/__TEMPLATE_SYSTEMD_BACKEND_STAGING__.service on the VM.
# Path: /opt/__TEMPLATE_DEPLOY_APP_NAME__/backend/staging/bin/backend
# Set APP_NAME via VM_APP_NAME secret (default: __TEMPLATE_DEPLOY_APP_NAME__).

[Unit]
Description=__TEMPLATE_APP_NAME__ Backend API (staging)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/__TEMPLATE_DEPLOY_APP_NAME__/backend/staging/bin/backend
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/__TEMPLATE_DEPLOY_APP_NAME__/backend/staging
EnvironmentFile=-/etc/default/__TEMPLATE_SYSTEMD_BACKEND_STAGING__

[Install]
WantedBy=multi-user.target
