# infra/systemd/backend-staging.service.example
# Example systemd unit for backend (staging).
# Copy to /etc/systemd/system/backend-staging.service on the VM.
# Path: /opt/<APP_NAME>/backend/staging/bin/backend
# Set APP_NAME via VM_APP_NAME secret (default: app).

[Unit]
Description=Backend API (staging)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/app/backend/staging/bin/backend
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/app/backend/staging
EnvironmentFile=-/etc/default/backend-staging

[Install]
WantedBy=multi-user.target
