# infra/systemd/backend-production.service.example
# Example systemd unit for backend (production).
# Copy to /etc/systemd/system/backend-production.service on the VM.
# Path: /opt/<APP_NAME>/backend/production/bin/backend

[Unit]
Description=Backend API (production)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/app/backend/production/bin/backend
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/app/backend/production
EnvironmentFile=-/etc/default/backend-production

[Install]
WantedBy=multi-user.target
