# infra/systemd/frontend-production.service.example
# Example systemd unit for frontend (production).
# Copy to /etc/systemd/system/frontend-production.service on the VM.

[Unit]
Description=Frontend (production)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/app/frontend/production
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/default/frontend-production

[Install]
WantedBy=multi-user.target
