# infra/systemd/frontend-staging.service.example
# Example systemd unit for frontend (staging).
# Copy to /etc/systemd/system/frontend-staging.service on the VM.
# Serves from /opt/<APP_NAME>/frontend/staging (e.g. SvelteKit adapter-node build).
# Adjust ExecStart if using a static server (nginx/caddy) instead of Node.

[Unit]
Description=Frontend (staging)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/app/frontend/staging
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/default/frontend-staging

[Install]
WantedBy=multi-user.target
